# build/authenticator/supervisord.conf
# Dropped into /etc/supervisor/conf.d/ ; the Debian-shipped
# /etc/supervisor/supervisord.conf provides the [supervisord]/rpc sections
# that fcgi-program needs and includes conf.d/*.conf.

[program:shibd]
command=/usr/sbin/shibd -F
user=_shibd
priority=10
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# fcgi-program makes supervisor open the listening socket and hand it to the
# binary on stdin -- which is exactly how the Shibboleth FastCGI apps expect
# to be launched.
[fcgi-program:shibauthorizer]
command=/usr/lib/x86_64-linux-gnu/shibboleth/shibauthorizer
socket=unix:///run/shibboleth/shibauthorizer.sock
socket_owner=_shibd:_shibd
socket_mode=0660
user=_shibd
priority=20
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[fcgi-program:shibresponder]
command=/usr/lib/x86_64-linux-gnu/shibboleth/shibresponder
socket=unix:///run/shibboleth/shibresponder.sock
socket_owner=_shibd:_shibd
socket_mode=0660
user=_shibd
priority=20
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=30
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
